# Server/Linux/service/php-fastcgi/run
#!/bin/sh
echo starting
exec sh /usr/bin/php-fastcgi-service

